# VitoryFairy_BE/py-collector/deploy/local/com.victoryfairy.fmkorea-catchup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd LaunchAgent: FMKorea popular-KBO top-up every 30 minutes.

  WHY launchd, not cron: cron skips runs while the Mac is asleep; launchd fires a
  missed StartInterval ONCE on the next wake. Combined with catchup_fmkorea.sh's
  lookback window, gaps up to LOOKBACK_DAYS are filled the next time the laptop is
  on. Idempotent (S3 dedup), so frequent runs only fetch posts not yet landed.
  launchd never starts a second instance while one is still running, and
  catchup_fmkorea.sh itself skips when a manual backfill_fmkorea.sh is active,
  so runs cannot pile up on FMKorea.

  INSTALL (FIRST replace every /ABS/PATH/TO/REPO below with your clone's location,
  e.g. /Users/you/VictoryFairy -- launchd needs absolute paths, no ~ or env expansion):
    cp deploy/local/com.victoryfairy.fmkorea-catchup.plist ~/Library/LaunchAgents/
    # then edit the copied file's 3 paths, and:
    launchctl load ~/Library/LaunchAgents/com.victoryfairy.fmkorea-catchup.plist
  RUN NOW (test):   launchctl start com.victoryfairy.fmkorea-catchup
  STATUS:           launchctl list | grep fmkorea-catchup
  UNINSTALL:        launchctl unload ~/Library/LaunchAgents/com.victoryfairy.fmkorea-catchup.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.victoryfairy.fmkorea-catchup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/ABS/PATH/TO/REPO/VitoryFairy_BE/py-collector/deploy/local/catchup_fmkorea.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- Target bucket for the dev stack. AWS creds come from ~/.aws automatically. -->
        <key>COLLECTOR_S3_BUCKET</key>
        <string>victoryfairy-crawl-dev</string>
        <!-- launchd starts with a bare PATH; add Homebrew (aws) + system dirs. -->
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <!-- 3-day window: each 30-min run re-walks only ~7-8 popularity-sorted list
             pages (S3-deduped), yet a weekend with the laptop off still gets filled
             on the next run. Widen if you skip more days. -->
        <key>LOOKBACK_DAYS</key>
        <string>3</string>
        <!-- One pass per firing: the next firing is 30 min away anyway, so the
             runner's multi-pass/block-cooldown loop would only add FMKorea load. -->
        <key>MAX_RUNS</key>
        <string>1</string>
    </dict>

    <!-- Every 30 minutes while awake; one catch-up firing after a sleep/boot gap. -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <!-- Don't fire the moment you `launchctl load` it; wait for the schedule. -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/ABS/PATH/TO/REPO/VitoryFairy_BE/py-collector/logs/launchd-catchup.out</string>
    <key>StandardErrorPath</key>
    <string>/ABS/PATH/TO/REPO/VitoryFairy_BE/py-collector/logs/launchd-catchup.err</string>
</dict>
</plist>
